# macos/com.user.dotfiles-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.dotfiles-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/temullen/dotfiles/bin/bin/dotfiles-sync</string>
    </array>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>StandardOutPath</key>
    <string>/Users/temullen/Library/Logs/dotfiles-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/temullen/Library/Logs/dotfiles-sync.log</string>
</dict>
</plist>
